# timer/teste-timer.service
[Unit]
Description=Executa teste-timer.sh
Wants=teste-timer.timer

[Service]
Type=oneshot
ExecStart=/usr/bin/free
ExecStart=/usr/bin/teste-timer.sh
ExecStart=/usr/bin/lsblk

[Install]
WantedBy=multi-user.target
